# dist/8.0/fpm-apache/etc/s6-overlay/s6-rc.d/apache2/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
/usr/sbin/apache2ctl -DFOREGROUND