# src/s6/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d
/usr/sbin/php-fpm --nodaemonize