# src/fpm/etc/s6-overlay/s6-rc.d/laravel-automations/up
/etc/s6-overlay/scripts/laravel-automations