# dist/8.1/fpm/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
with-contenv
exec /usr/sbin/php-fpm8.1 --nodaemonize