# generated-dockerfiles/8.1/fpm-apache/etc/services.d/apache2/run
#!/usr/bin/with-contenv sh
/usr/sbin/apache2ctl -DFOREGROUND