# dist/8.1/fpm-nginx/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
with-contenv
nginx