# src/fpm-apache/etc/services.d/apache2/run
#!/command/with-contenv sh
/usr/sbin/apache2ctl -DFOREGROUND