# dist/8.1/fpm-apache/etc/services.d/apache2/run
#!/command/with-contenv bash
/usr/sbin/apache2ctl -DFOREGROUND