# src/fpm-apache/etc/s6-overlay/s6-rc.d/apache2/run
#!/command/execlineb -P
with-contenv

# Set healthcheck sleep delay to 20ms, because you know... Apache
s6-notifyoncheck -s 20
/usr/sbin/apache2ctl -DFOREGROUND