# php/7.4/fpm/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh
su webuser
exec /usr/sbin/php-fpm7.4 --nodaemonize