# php/7.4/fpm/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

# Switch to web user
su webuser

#Bring up PHP FPM
exec /usr/sbin/php-fpm7.4 --nodaemonize