# dist/8.1/fpm-apache/etc/s6-overlay/s6-rc.d/apache2/run
#!/command/execlineb -P
with-contenv
/usr/sbin/apache2ctl -DFOREGROUND