# dist/7.4/fpm/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
exec /usr/sbin/php-fpm7.4 --nodaemonize